# board/chumby/chumby8/rootfs-overlay/etc/profile.d/display.sh
# Export the DISPLAY variable so we can run X11 apps from the command line.
export DISPLAY=:0
